# templates/docker-prune.service.j2
[Unit]
Description=Prune unused Docker images

[Service]
Type=simple
ExecStart=/usr/bin/docker image prune --force

[Install]
WantedBy=default.target